# transparente-ugr.conf
#!upstart
description "transparente-ugr"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec nodejs /var/www/transparente-ugr/app.js >> /var/log/transparente-ugr.log 2>&1