# provisioner/transparente-ugr.conf
start on startup
stop on shutdown

env PORT=80
exec forever start -l /var/log/forever.log -a -o /var/log/out.log -e /var/log/err.log /home/{user}/ugr-transparente-servidor/app.js
